# start-clean-none-docker-image.service
[Unit]
Description=Run clean none tag docker images every day at 0 AM

[Timer]
OnCalendar=00:00
Persistent=true

[Install]
WantedBy=local.target
